# watcher.service
[Unit]
Description=Python Connection Watcher Bot
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /path/to/main.py
Restart=always
User=your_username
WorkingDirectory=/path/to/bot/directory
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
